# playbooks/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[program:sshd]
command=/usr/sbin/sshd -D

