# playbooks/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[program:mysqld]
user=docker
command=/usr/libexec/mysqld

[program:sshd]
command=/usr/sbin/sshd -D
